# agent/inband/config/vlan.config
### BEGIN INIT INFO
# Description: This file utilized by olt inband init script and onl validation script to
#              know what VLAN configured by the user for inband CTRL and MGMT interface
### END INIT INFO

#asfvolt16
asfvolt16_vlan_id_eth2=4093

#asgvolt64
asgvolt64_vlan_id_eth1=4093
